# systemd/vmquota-api.service
[Unit]
Description=vmquota read-only API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/vmquota serve
Restart=always
RestartSec=2s

[Install]
WantedBy=multi-user.target
